# homelab-dashboard.service
[Unit]
Description=Homelab Dashboard - UniFi and Kubernetes Management
Documentation=https://github.com/imaustink/admin-console
After=network.target graphical-session.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/opt/homelab-dashboard

# Display (required for Tauri/WebKit)
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/admin/.Xauthority"

# Disable problematic WebKit GPU rendering on Raspberry Pi
Environment="WEBKIT_DISABLE_DMABUF_RENDERER=1"
Environment="WEBKIT_DISABLE_COMPOSITING_MODE=1"

# Config path
Environment="CONFIG_PATH=/etc/homelab-dashboard/config.json"

ExecStart=/opt/homelab-dashboard/homelab-dashboard.AppImage --no-sandbox

Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=homelab-dashboard

[Install]
WantedBy=graphical-session.target
